# deploy/petshare.service
[Unit]
Description=PetShare Israel Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/petshare
ExecStart=/opt/petshare/venv/bin/python bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
